# examples/gateway.service
[Unit]
Description=Wack-A-Mole Gateway
After=network.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/gateway.env
WorkingDirectory=/opt/wack-gateway
ExecStart=/usr/bin/uvicorn gateway:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=3
LimitNOFILE=65536

# Hardening (adjust if needed)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target